# conf/supervisor/worker-supervisord.conf


# File: worker-supervisor.conf
# Created at 12/04/2023

[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
loglevel=info                ; (log level;default info; others: debug,warn,trace)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL  for a unix socket



[program:consumer-service-1]
command=celery --app worker worker -Q guard-evident-1-queue -l INFO -c 4
directory=/webapps/service
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:consumer-service-2]
command=celery --app worker worker -Q guard-evident-2-queue -l INFO -c 4
directory=/webapps/service
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:consumer-service-3]
command=celery --app worker worker -Q guard-evident-3-queue -l INFO -c 4
directory=/webapps/service
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:consumer-service-4]
command=celery --app worker worker -Q guard-evident-4-queue -l INFO -c 4
directory=/webapps/service
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:consumer-service-5]
command=celery --app worker worker -Q guard-evident-5-queue -l INFO -c 4
directory=/webapps/service
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0